# scripts/jon-server/profile.d/jon-server.sh
#!/usr/bin/env bash
export RHQ_JAVA_HOME="/usr/lib/jvm/jre"
export RHQ_SERVER_HOME="/apps/jon/jon-server-3.3.0.GA"
export RHQ_SERVER_USER="jon"
